# crontab
@reboot sleep 15 && /root/startup 2>&1
* * * * * /root/servicecheck 2>&1

