# crontab
@reboot /usr/bin/sh /root/startup.sh
* * * * * /usr/bin/sh /root/servicecheck.sh

