# crontab
@reboot sleep 15 && /bin/sh /root/startup
* * * * * /bin/sh /root/servicecheck

