# crontab
@reboot sleep 15 && /root/startup
* * * * * /root/servicecheck

