# crontab
@reboot /usr/bin/sh /root/startup
* * * * * /root/servicecheck

